# quadlets/openclaw.target
[Unit]
Description=OpenClaw Zoo

[Install]
WantedBy=default.target
