# quadlets/openclaw-cli.container
[Unit]
Description=OpenClaw CLI
After=network-online.target
Wants=network-online.target

[Container]
Image=${OPENCLAW_IMAGE:-openclaw:local}
Volume=openclaw-config:/home/node/.openclaw
Volume=openclaw-workspace:/home/node/.openclaw/workspace
Environment=HOME=/home/node
Environment=TERM=xterm-256color
Environment=BROWSER=echo
Environment=OPENCLAW_GATEWAY_TOKEN=${OPENCLAW_GATEWAY_TOKEN}
Environment=CLAUDE_AI_SESSION_KEY=${CLAUDE_AI_SESSION_KEY}
Environment=CLAUDE_WEB_SESSION_KEY=${CLAUDE_WEB_SESSION_KEY}
Environment=CLAUDE_WEB_COOKIE=${CLAUDE_WEB_COOKIE}
Exec=node dist/index.js

[Install]
WantedBy=multi-user.target
